# vpn-container/root/etc/s6-overlay/s6-rc.d/svc-dnsbl-checker/run
#!/usr/bin/with-contenv bash
exec /usr/local/bin/dnsbl-checker 2>&1
# vi: set ft=bash:
